# templates/com.waconcierge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- wa-concierge - launchd agent (macOS).
     The installer copies this to ~/Library/LaunchAgents/, replacing
     NODE_PATH_AQUI (output of `which node`) and PROJETO_AQUI (absolute project path).
     Load with: launchctl load ~/Library/LaunchAgents/com.waconcierge.plist -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.waconcierge.server</string>
  <key>ProgramArguments</key>
  <array>
    <string>NODE_PATH_AQUI</string>
    <string>PROJETO_AQUI/src/web-entry.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>PROJETO_AQUI</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>PROJETO_AQUI/data/server.log</string>
  <key>StandardErrorPath</key>
  <string>PROJETO_AQUI/data/server.log</string>
</dict>
</plist>
